# files/etc/uci-defaults/99-jdc01-sovereignty
#!/bin/sh

# --- 1. RPS/RFS 核心对齐 ---
# f = 1111 (4个逻辑核心全开)
for rps_file in /sys/class/net/*/queues/rx-*/rps_cpus; do
    echo "f" > "$rps_file" 2>/dev/null
done

echo "2048" > /proc/sys/net/core/rps_sock_flow_entries
for rfc_file in /sys/class/net/*/queues/rx-*/rps_flow_cnt; do
    echo "512" > "$rfc_file" 2>/dev/null
done

# --- 2. 硬件加速与网卡特性 ---
for eth in $(ls /sys/class/net | grep -E 'eth|wan|lan'); do
    ethtool -K "$eth" gro on 2>/dev/null
done

# --- 3. 内存管理：zram + zstd ---
if [ -b /dev/zram0 ]; then
    swapoff /dev/zram0 2>/dev/null
    echo 1 > /sys/block/zram0/reset 2>/dev/null
    # 强制检查并设置 zstd
    [ -f /sys/block/zram0/comp_algorithm ] && echo zstd > /sys/block/zram0/comp_algorithm
    TOTAL_MEM=$(free | grep Mem | awk '{print $2}')
    # 分配 50% 内存作为压缩 Swap
    echo "$((TOTAL_MEM / 2 * 1024))" > /sys/block/zram0/disksize
    mkswap /dev/zram0 >/dev/null
    swapon -p 10 /dev/zram0
fi

# --- 4. MT7621 养老与 TCP 调优 ---
cat <<EOF >> /etc/sysctl.conf
net.netfilter.nf_conntrack_max=65535
net.netfilter.nf_conntrack_tcp_timeout_established=600
net.netfilter.nf_conntrack_tcp_timeout_close_wait=30
net.core.netdev_max_backlog=2048
net.ipv4.tcp_fastopen=3
net.core.default_qdisc=fq_codel
net.ipv4.tcp_congestion_control=bbr
vm.swappiness=100
vfs_cache_pressure=10
dirty_ratio=10
dirty_background_ratio=5
vm.min_free_kbytes = 4096
EOF
sysctl -p

#!/bin/sh

# --- 5.Kernel极度优化调参
# 1. 基础逻辑优化
uci set network.globals.packet_steering='1'
uci set firewall.@defaults[0].flow_offloading='0'
uci set firewall.@defaults[0].flow_offloading_hw='0'
uci commit firewall

# 2. 针对 JDC01 的 RPS 深度绑定
# JDC01 通常只有一个 MT7615，我们将中断处理尽量分散到 CPU1/2/3
logger "Applying RPS tuning for JDC01 single-chip architecture..."
for dev in eth0 ra0 rai0 br-lan; do
    [ -d "/sys/class/net/$dev" ] && {
        # 允许所有核心参与处理，减少 CPU0 的软中断压力
        echo f > /sys/class/net/$dev/queues/rx-0/rps_cpus
        echo 1024 > /sys/class/net/$dev/queues/rx-0/rps_flow_cnt
    }
done

# 3. 增强型网络内核调优 (针对大流量 4K 压强)
cat <<EOF > /etc/sysctl.d/99-network-pressure.conf
net.core.rmem_max=16777216
net.core.wmem_max=16777216
net.ipv4.tcp_rmem=4096 87380 16777216
net.ipv4.tcp_wmem=4096 65536 16777216
net.core.somaxconn=4096
EOF

# 4. 隧道气流优化 (针对翻墙应用)
cat <<EOF > /etc/hotplug.d/iface/99-jdc-tunnel
[ "\$ACTION" = "ifup" ] && {
    # 强制让隧道接口利用非 0 核心，避开无线中断主核
    echo e > /sys/class/net/\$INTERFACE/queues/rx-0/rps_cpus 2>/dev/null
}
EOF

# 5. 硬件层 Ring Buffer 扩容 (ethtool 注入)
# 直接将以太网物理接口的收发描述符队列拉满到 2048
if command -v ethtool >/dev/null 2>&1; then
    ethtool -G eth0 rx 2048 tx 2048 >/dev/null 2>&1
fi

# 6. 链路聚合哈希与接口发送队列调优
# 强制 L3+L4 哈希策略
if [ -d /sys/class/net/bond0 ]; then
    echo "layer3+4" > /sys/class/net/bond0/bonding/xmit_hash_policy
fi

# 7. 中断合并调优 (Interrupt Coalescing)
# 让网卡每积攒 16 个包或者每过 50 微秒才去“烦”一次 CPU
if command -v ethtool >/dev/null 2>&1; then
    ethtool -C eth0 rx-frames 16 rx-usecs 50 tx-frames 16 tx-usecs 50 >/dev/null 2>&1
fi
# 增加接口发送队列长度，配合 PPE 降低软中断负载
for dev in ra0 rai0 bond0 eth0; do
    if ifconfig $dev > /dev/null 2>&1; then
        ifconfig $dev txqueuelen 2048
    fi
done

# --- 6. 无线环境极致固化 ---
for radio in $(uci show wireless | grep =wifi-device | cut -d. -f2 | cut -d= -f1); do
    uci set wireless.${radio}.noscan='1'
    uci set wireless.${radio}.country='AU'
    uci set wireless.${radio}.legacy_rates='0'
    uci set wireless.${radio0}.bgscan='0' 
done

logger -t "MIPS-Optimizer" "MT7621 内核对齐完成：RPS(掩码F)已衔接，MRS 规则已对位。"
exit 0
