# files/etc/uci-defaults/99-physical-sovereignty
#!/bin/sh
# --- MT7621 + MT7615D 物理主權 v2.2 (極限煉丹版) ---

sleep 10

# 1. 物理線程綁定 (IRQ Affinity)
# 將無線中斷鎖定在 CPU2/3 (物理第二核)，有線鎖定在 CPU0
ETH_IRQ=$(grep -m1 "ethernet" /proc/interrupts | cut -d: -f1 | tr -d ' ')
WIFI_IRQS=$(grep 'mt76' /proc/interrupts | cut -d: -f1)

[ -n "$ETH_IRQ" ] && echo 1 > "/proc/irq/$ETH_IRQ/smp_affinity"
for irq in $WIFI_IRQS; do
    echo 4 > "/proc/irq/$irq/smp_affinity"
done

# 2. 徹底關閉 RPS (去噪，回歸物理核心響應)
for rps in /sys/class/net/*/queues/rx-*/rps_cpus; do echo 0 > "$rps"; done

# 3. 內存雙引擎：KSM + zRAM (256M 物理內存折疊)
# [KSM] 溫和掃描，不搶 CPU 週期
if [ -f /sys/kernel/mm/ksm/run ]; then
    echo 50 > /sys/kernel/mm/ksm/pages_to_scan
    echo 200 > /sys/kernel/mm/ksm/sleep_millisecs
    echo 1 > /sys/kernel/mm/ksm/run
fi

# [zRAM] 針對 1.1GHz 超頻優化：zstd 高壓縮
if [ -d /sys/block/zram0 ]; then
    swapoff /dev/zram0 2>/dev/null
    echo zstd > /sys/block/zram0/comp_algorithm
    echo 128M > /sys/block/zram0/disksize
    mkswap /dev/zram0 >/dev/null
    swapon -p 10 /dev/zram0
fi

# 4. 網卡卸載與中斷聚合 (1.1GHz 冷靜補丁)
ethtool -K eth0 gro on gso on tso on 2>/dev/null
ethtool -C eth0 rx-usecs 150 tx-usecs 150 2>/dev/null

# 5. MT7615D 無線私有優化 (中繼核心補丁)
# 鎖死中繼掃描頻率，開啟 VHT 加速
uci batch <<EOF
# 關閉 5G 背景掃描，防止 CS/吃雞 跳 Ping
set wireless.radio0.bgscan='0'
set wireless.radio1.bgscan='0'
# 設置 VHT160 或 80 高性能模式 (視驅動支持)
set wireless.radio1.vht_mimo='1'
set wireless.radio1.txqueuelen='5000'
# 禁用無線省電 (物理級鎖定)
set wireless.default_radio1.disassoc_low_ack='0'
set wireless.radio1.noscan='1'
# 強制開啟 HNAT (如果底層支持)
set network.lan.hwnat='1'
EOF
uci commit wireless
uci commit network

# 6. HNAT 中繼補丁 (物理喚醒)
if [ -d "/sys/kernel/debug/hnat" ]; then
    echo 1 > /sys/kernel/debug/hnat/hook_relay
    echo 1 > /sys/kernel/debug/hnat/all_entry_keepalive
fi

# 7. 內核消噪與精煉 TCP 棧 (GFW 屏蔽方案)
cat >> /etc/sysctl.conf <<EOF
net.core.default_qdisc=fq_codel
net.ipv4.tcp_congestion_control=bbr
net.core.rmem_max=2097152
net.core.wmem_max=2097152
net.ipv4.tcp_fastopen=3
net.ipv4.tcp_fin_timeout=15
net.ipv4.icmp_echo_ignore_broadcasts=1
vm.swappiness=60
vm.min_free_kbytes=4096
EOF

# 8. 遊戲 UDP 染色優先 (NFTables)
nft insert rule inet fw4 mangle_postrouting udp dport { 1024-65535 } udp length 0-128 ip dscp set ef 2>/dev/null

sysctl -p
exit 0
